# lib/systemd/clamd.service.sample
[Unit]
Description=ClamAV Antivirus Daemon (clamd)
Before=multi-user.target

[Service]
ExecStart=/usr/local/sbin/clamd
Type=forking

[Install]
WantedBy=amavis.service
WantedBy=multi-user.target
