# lib/systemd/psad.service.sample
[Unit]
Description=Electronic Log Daemon (elogd)
After=multiuser.target

[Service]
ExecStart=/usr/sbin/psad
Type=forking

[Install]
WantedBy=multi-user.target
Requires=iptables.service
